# systemd/arclap-station.service
[Unit]
Description=Arclap Station — FastAPI control plane
Documentation=https://github.com/arclap-af/arclap-station
After=network-online.target arclap-station.socket
Wants=network-online.target
Requires=arclap-station.socket

[Service]
Type=notify
User=arclap
Group=arclap
WorkingDirectory=/var/lib/arclap
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=ARCLAP_CONF_DIR=/etc/arclap
Environment=ARCLAP_STATE_DIR=/var/lib/arclap
Environment=ARCLAP_PHOTO_DIR=/media/sdcard/photos
Environment=ARCLAP_LOG_DIR=/var/log/arclap
EnvironmentFile=-/etc/arclap/environment

# systemd passes FD 3 as the listening socket. arclap-station's CLI reads it.
ExecStart=/opt/arclap-station/venv/bin/arclap-station serve --socket-fd=3
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=3
TimeoutStartSec=60
TimeoutStopSec=30

# Hardware watchdog — the Pi 5 BCM2712 SoC has /dev/watchdog.
# 60s grace; the app emits WATCHDOG=1 over $NOTIFY_SOCKET every 15s.
RuntimeWatchdogSec=60s
WatchdogSec=30s

# Resource limits.
LimitNOFILE=4096
TasksMax=512
MemoryMax=512M

# ---- Hardening ----
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
DeviceAllow=char-usb_device rw
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/bus/usb rw
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Bind a privileged port? No — we listen on a UDS, Caddy fronts 443.
CapabilityBoundingSet=
AmbientCapabilities=

# Writable paths.
ReadWritePaths=/etc/arclap /var/lib/arclap /var/log/arclap /media/sdcard
ReadOnlyPaths=/opt/arclap-station

# Standard streams → journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arclap-station

[Install]
WantedBy=multi-user.target
Also=arclap-station.socket
