# systemd/arclap-station.service
[Unit]
Description=Arclap Station — FastAPI control plane
Documentation=https://github.com/arclap-af/arclap-station
After=network-online.target
Wants=network-online.target

[Service]
# Plain TCP — uvicorn binds 127.0.0.1:8080, Caddy reverse-proxies to that.
# Socket activation was planned but the CLI doesn't implement the
# LISTEN_FDS protocol or sd_notify, so Type=simple is the honest choice.
Type=simple
User=arclap
Group=arclap
WorkingDirectory=/var/lib/arclap
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=ARCLAP_CONF_DIR=/etc/arclap
Environment=ARCLAP_STATE_DIR=/var/lib/arclap
Environment=ARCLAP_PHOTO_DIR=/media/sdcard/photos
Environment=ARCLAP_LOG_DIR=/var/log/arclap
Environment=ARCLAP_BIND_HOST=127.0.0.1
Environment=ARCLAP_BIND_PORT=8080
EnvironmentFile=-/etc/arclap/environment

ExecStart=/opt/arclap-station/venv/bin/arclap-station serve
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=3
TimeoutStartSec=60
TimeoutStopSec=30

# Resource limits.
LimitNOFILE=4096
TasksMax=512
MemoryMax=512M

# ---- Hardening ----
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
DeviceAllow=char-usb_device rw
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/bus/usb rw
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# We bind 127.0.0.1:8080 (unprivileged) — Caddy fronts 443.
CapabilityBoundingSet=
AmbientCapabilities=

# Writable paths.
ReadWritePaths=/etc/arclap /var/lib/arclap /var/log/arclap /media/sdcard
ReadOnlyPaths=/opt/arclap-station

# Standard streams → journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arclap-station

[Install]
WantedBy=multi-user.target
