# systemd/arclap-timelapse.timer
[Unit]
Description=Daily 24h timelapse render
Documentation=https://github.com/arclap-af/arclap-station

[Timer]
# 03:30 local — after retention sweep (03:00), before backup (04:00).
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300
AccuracySec=1min

[Install]
WantedBy=timers.target
