# packaging/deb/postinst
#!/bin/sh
# Arclap Station — .deb postinst hook.
# Runs after dpkg has dropped the payload into place.

set -e

ARCLAP_USER=arclap
ARCLAP_HOME=/var/lib/arclap
ARCLAP_PREFIX=/opt/arclap-station
ARCLAP_CONFDIR=/etc/arclap
ARCLAP_LOGDIR=/var/log/arclap
ARCLAP_PHOTODIR=/media/sdcard/photos

case "$1" in
  configure)
    # 1. Create system user (idempotent).
    if ! id "$ARCLAP_USER" >/dev/null 2>&1; then
      adduser --system --quiet --group --home "$ARCLAP_HOME" --no-create-home \
        --shell /usr/sbin/nologin "$ARCLAP_USER"
    fi
    for g in plugdev video dialout; do
      if getent group "$g" >/dev/null 2>&1; then
        usermod -a -G "$g" "$ARCLAP_USER" 2>/dev/null || true
      fi
    done

    # 2. State dirs.
    install -d -m 0750 -o "$ARCLAP_USER" -g "$ARCLAP_USER" \
      "$ARCLAP_HOME" "$ARCLAP_CONFDIR" "$ARCLAP_CONFDIR/destinations" \
      "$ARCLAP_LOGDIR" "$ARCLAP_PHOTODIR"

    # 3. Install the wheel into a venv. Pick a python3.x that exists
    # (Pi-OS Bookworm ships 3.11, Ubuntu Resolute ships 3.14, etc.).
    ARCLAP_PY=""
    for candidate in python3.13 python3.12 python3.11 python3; do
      if command -v "$candidate" >/dev/null 2>&1; then
        ARCLAP_PY="$candidate"
        break
      fi
    done
    if [ -z "$ARCLAP_PY" ]; then
      echo "arclap-station postinst: no python3 found in PATH; aborting" >&2
      exit 1
    fi
    if [ ! -x "$ARCLAP_PREFIX/venv/bin/arclap-station" ]; then
      "$ARCLAP_PY" -m venv "$ARCLAP_PREFIX/venv"
      "$ARCLAP_PREFIX/venv/bin/pip" install --quiet --upgrade pip
      whl=$(ls "$ARCLAP_PREFIX"/wheels/*.whl 2>/dev/null | head -n1)
      if [ -n "$whl" ]; then
        "$ARCLAP_PREFIX/venv/bin/pip" install --quiet "$whl"
      fi
    fi
    ln -sfn "$ARCLAP_PREFIX/venv/bin/arclap-station" /usr/local/bin/arclap-station

    chown -R "$ARCLAP_USER:$ARCLAP_USER" "$ARCLAP_PREFIX"

    # 4. Generate Caddyfile from template.
    serial=$(awk -F: '/^Serial/ { gsub(/ /, "", $2); print tolower($2); exit }' /proc/cpuinfo 2>/dev/null \
             | tail -c 9 | tr -d '\n')
    if [ -z "$serial" ]; then serial="unknown1"; fi
    if [ -f /etc/caddy/Caddyfile.template ]; then
      sed "s|\${SERIAL}|$serial|g" /etc/caddy/Caddyfile.template > /etc/caddy/Caddyfile
    fi

    # 5. Hostname. Accept the common defaults across Pi-OS + Ubuntu
    # (raspberrypi, arclap-st, ubuntu) so first-time installs always rename.
    cur_host="$(hostnamectl --static 2>/dev/null || echo "")"
    case "$cur_host" in
      raspberrypi|arclap-st|ubuntu|"")
        hostnamectl set-hostname "arclap-st-$serial" || true
        grep -q "arclap-st-$serial" /etc/hosts || \
          printf "127.0.1.1\t%s\n" "arclap-st-$serial" >> /etc/hosts
        ;;
    esac

    # 6. Reload + enable. Socket / uploader units are NOT enabled — the
    # CLI doesn't implement socket activation and the uploader runs
    # in-process via FastAPI lifespan.
    udevadm control --reload-rules || true
    udevadm trigger --action=change || true
    systemctl daemon-reload
    deb-systemd-helper enable arclap-station.service >/dev/null 2>&1 || true
    deb-systemd-helper enable arclap-watchdog.timer >/dev/null 2>&1 || true

    if [ -d /run/systemd/system ]; then
      systemctl restart caddy avahi-daemon 2>/dev/null || true
      systemctl restart arclap-station.service 2>/dev/null || true
    fi
    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;
esac

exit 0
